# systemd/smartclock.service
[Unit]
Description=SmartClock

[Service]
Type=forking
User=root
ExecStart=/home/pi/smartclock/smartclock_daemon.sh

[Install]
WantedBy=multi-user.target